# etc/cron.d/bx_dbtest-vallesru
#
# cron tasks for site /home/bitrix/ext_www/test-vallesru.shellbe.ru db=dbtest-vallesru
#
* * * * * bitrix test -f /home/bitrix/ext_www/test-vallesru.shellbe.ru/bitrix/modules/main/tools/cron_events.php && { /usr/bin/php -d sendmail_path="msmtp -t -i -a test-vallesru.shellbe.ru" -f /home/bitrix/ext_www/test-vallesru.shellbe.ru/bitrix/modules/main/tools/cron_events.php; } >/dev/null 2>&1
