# etc/cron.d/bx_httpd-scale
* * * * * root /opt/webdir/bin/restart_httpd-scale.sh process
