# etc/cron.d/bx_network_updater
# Ansible managed
# test nginx log for network update requests
*/15 * * * * root /opt/webdir/bin/wrapper_ansible_conf -a check_network --update_log /opt/webdir/logs/update_pool.log 1>/dev/null 2>&1
